# profile.d/xdg.sh
# Variables generales:
export XDG_CONFIG_HOME="~/.config"
export XDG_DATA_HOME="$XDG_CONFIG_HOME/"
export XDG_CACHE_HOME="$XDG_CONFIG_HOME/.cache"